# confluence.service.j2
[Unit]
Description=Confluence
After=network.target

[Service]
Type=forking
User=confluence
PIDFile={{ confluence_install_dir }}/work/catalina.pid
ExecStart={{ confluence_install_dir }}/bin/start-con
