# jira.service.j2
[Unit]
Description=JIRA
After=network.target

[Service]
Type=forking
User=jira
PIDFile={{ jira_install_dir }}/work/catalina.pid
ExecStart
